# deployment/erp-system.service
[Unit]
Description=ERP System Next.js application
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=erp
Group=erp
WorkingDirectory=/var/www/erp/current
EnvironmentFile=/var/www/erp/current/.env
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=PATH=/usr/bin:/usr/local/bin
Environment=HOME=/var/www/erp
ExecStart=/usr/bin/npm run start:prod
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
LimitNOFILE=65536
SyslogIdentifier=erp-system
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
